# WDMyCloud Firmware/firmware/module/crfs/etc/cron.d/rescan_internal_volumes
10 3 * * * root  [ -f /tmp/wdmc_rescan ] && /usr/local/sbin/wdmc_rescan_volume.py
